# .config/autostart/aw-qt.desktop
[Desktop Entry]
Type=Application
Name=Activity Watch
Exec=aw-qt
StartupNotify=false
Terminal=false
Categories=Monitor;
